# rustic/rustic-backup.timer
[Unit]
Description=Backup every 3 hours

[Timer]
OnStartupSec=20min
OnUnitInactiveSec=3h

[Install]
WantedBy=timers.target
