# rustic/rustic-backup.service
[Unit]
Description=Backup with rustic

[Service]
ExecStart=/usr/bin/rustic backup
